# deploy/systemd/airwv-reference.timer
# Runs airwv-reference.service every 4 hours. OpenAQ's free tier has a DAILY
# request quota (exceeding it gets the key suspended), so we keep ongoing volume
# well under it: ~monitors×6/day, plus the hard OPENAQ_DAILY_CAP in the app as a
# backstop. EPA/AirNow data lands late anyway, so 4h + `--days 2` loses nothing.
#   sudo systemctl enable --now airwv-reference.timer
# Check: systemctl list-timers airwv-reference.timer

[Unit]
Description=Run AirWV reference-monitor collection every 4h (within OpenAQ daily quota)

[Timer]
OnBootSec=8min
OnUnitActiveSec=4h
Persistent=true

[Install]
WantedBy=timers.target
