# deploy/systemd/airwv-collect.service
# AirWV collection — a one-shot run triggered by airwv-collect.timer.
# Copy both files to /etc/systemd/system/, edit paths/user, then:
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now airwv-collect.timer
#
# Adjust WorkingDirectory, User, and the venv path to your install.

[Unit]
Description=AirWV — collect PurpleAir readings into storage
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=airwv
WorkingDirectory=/opt/airwv
EnvironmentFile=/opt/airwv/.env
ExecStart=/opt/airwv/.venv/bin/python -m airwv.ingest collect
