# deploy/systemd/airwv-collect.timer
# Runs airwv-collect.service hourly (PurpleAir points add up per sensor per call;
# hourly is sustainable on the free tier for ~30 sensors). Enable with:
#   sudo systemctl enable --now airwv-collect.timer
# Check status: systemctl list-timers airwv-collect.timer

[Unit]
Description=Run AirWV collection hourly

[Timer]
OnBootSec=5min
OnUnitActiveSec=6h
Persistent=true

[Install]
WantedBy=timers.target
