# deploy/systemd/airwv-water.service
# USGS real-time water collection — one-shot, triggered by airwv-water.timer.
# Pulls USGS NWIS instantaneous water readings for WV (KEYLESS — no quota): pH,
# dissolved oxygen, conductivity, turbidity, temperature, streamflow, gage height.
# Copy both to /etc/systemd/system/, then:
#   sudo systemctl daemon-reload && sudo systemctl enable --now airwv-water.timer

[Unit]
Description=AirWV — collect USGS real-time water readings for WV (keyless)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=airwv
WorkingDirectory=/opt/airwv
EnvironmentFile=/opt/airwv/.env
ExecStart=/opt/airwv/.venv/bin/python -m airwv.ingest water
